# cover.service
[Unit]
Description=COVER

[Service]
ExecStart=/home/root/cover.arm

[Install]
WantedBy=multi-user.target
